# etc/consul.upstart.conf
description "Consul server process"

start on (local-filesystems and net-device-up IFACE=eth1)
stop on runlevel [!12345]

respawn

setuid consul
setgid consul

exec consul agent -config-dir /etc/consul.d/server -bind `ifconfig eth1|grep 'inet addr'|cut -d: -f2|awk '{print $1}'`
